# docscleaner_1.0-jammy_all/etc/systemd/system/docscleaner.service
[Unit]
Description=Downloads cleaner
After=network.target
StartLimitIntervalSec=0
[Service]
Type=simple
Restart=always
RestartSec=1
User=samuel
ExecStart=/usr/local/bin/docscli /home/$USER/Downloads

[Install]
WantedBy=multi-user.target
